# Module_0_Tasks/Task2_Backup_service/backup.service
[Unit]
Description=Automated backup
After=network.target

[Service]
ExecStart= /bin/backup.sh 


[Install]
WantedBy=multi-user.target

